# systemd/ndictd.service
[Unit]
Description=ndict - Speech to Text Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sh -c 'mkdir -p ~/.local/share/ndict'
ExecStart=/bin/sh -c '/home/uneeb/.local/bin/ndictd'
ExecStop=/bin/sh -c 'pkill -f ndictd; rm -f $XDG_RUNTIME_DIR/ndictd.sock 2>/dev/null || rm -f /tmp/ndictd.sock'
Restart=on-failure
RestartSec=5
AmbientCapabilities=CAP_SYS_INPUT
PassEnvironment=WAYLAND_DISPLAY DISPLAY XDG_SESSION_TYPE XDG_CURRENT_DESKTOP

[Install]
WantedBy=default.target
