# systemd/ndictd.service
[Unit]
Description=ndict - Speech to Text Daemon
# Wait for sound, but we handle the Wayland wait in the Service section
After=pipewire.service pulseaudio.service

[Service]
Type=simple

# 1. HARDCODED WAYLAND-1 ENVIRONMENT
Environment=WAYLAND_DISPLAY=wayland-1
Environment=DISPLAY=:1
Environment=XDG_RUNTIME_DIR=/run/user/%U

# 2. SMART WAIT: This loop pauses the service until Niri creates the wayland-1 socket.
# It checks every 0.5 seconds and won't proceed until the display is ready.
ExecStartPre=/bin/sh -c 'while [ ! -S "$XDG_RUNTIME_DIR/wayland-1" ]; do sleep 0.5; done'

# 3. CLEANUP: Ensure folders exist and old sockets are gone
ExecStartPre=/bin/sh -c 'mkdir -p %h/.local/share/ndict; rm -f $XDG_RUNTIME_DIR/ndictd.sock /tmp/ndictd.sock'

# 4. LAUNCH
ExecStart=%h/.local/bin/ndictd
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
